# deploy/vps/systemd/justinstone-comments.service
[Unit]
Description=Justin Stone first-party comments service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
EnvironmentFile=/etc/default/justinstone-comments
ExecStart=/usr/bin/python3 /srv/apps/justinstone.online/comments/comment_service.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/data/justinstone.online

[Install]
WantedBy=multi-user.target
